# server/lumora-server.service
[Unit]
Description=Lumora Local Hub Server
Documentation=https://github.com/lincon07/lumora_kiosk
After=network.target
Wants=network-online.target

[Service]
Type=simple

# Run as the kiosk user (replace "lumora" with your actual Linux username)
User=%i
Group=%i

# Working directory — server package root
WorkingDirectory=/opt/lumora/server

# Node.js binary (adjust path if using nvm or a different install location)
ExecStart=/usr/bin/node /opt/lumora/server/dist/index.js

# Restart on failure with a short backoff
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=5

# Environment
Environment=NODE_ENV=production
Environment=LUMORA_PORT=4000

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=%h/.lumora

# Logging — captured by journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumora-server

[Install]
WantedBy=multi-user.target
